# pkg/vmimage/data/etc/rc.d/rc.local
#!/bin/bash -ex

# https://github.com/openshift/openshift-azure/issues/1454
for pos in /var/log/journald.pos/worker0/storage.json /var/log/openshift-audit.pos; do
    if [[ -e $pos && ! -s $pos ]]; then
        rm $pos
    fi
done
